# launchd/com.nanoclaw.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>~/Dev/fullstackOS/services/nanoclaw/nanoclaw.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>~/Dev/fullstackOS/services/nanoclaw</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>~/.nanoclaw/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>~/.nanoclaw/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>~/Dev/fullstackOS/services/nanoclaw</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
